# supervisord.conf
[supervisord]
nodaemon=true

[program:terraria]
command=/usr/bin/mono /TerrariaServer.exe -configpath "/config" -config "/config/serverconfig.txt" -worldpath /world -logpath "/logs"
